# data/supervisord.conf
[inet_http_server]
port = *:9001

[unix_http_server]
file = /var/run/supervisor.sock
chmod = 0770
chown = nobody:nogroup

[supervisord]
nodaemon=true
logfile = /var/log/supervisor/supervisord.log
pidfile = /var/run/supervisord.pid

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[include]
files = /etc/supervisor/conf.d/*.conf

[program:cron]
command = /usr/sbin/cron -f
autostart = true
autorestart = true

[program:server]
command = php -S ${LISTEN_IP:-0.0.0.0}:${LISTEN_PORT:-8080} -t ${PROJECT_PATH:-/var/www/app}/public ${PROJECT_PATH:-/var/www/app}/public/index.php
autostart = true
autorestart = true